# scripts/battery-alert/.config/systemd/user/battery-alert.service
[Unit]
Description=Battery Alert Notification Service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/battery_alert.py
